# rootfs/usr/lib/systemd/system/handbus.service
[Unit]
Description=HanDBus Service
After=graphical-session.target

[Service]
ExecStart=/usr/bin/handbus

[Install]
WantedBy=multi-user.target
